# contrib/supervisor/redis-sentinel.conf
[program:redis-sentinel]
command=/usr/bin/redis-sentinel /etc/redis/sentinel.conf
autostart=true
autorestart=true
user=redis
